# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:laravel-queue-worker]
command=php /var/www/artisan queue:work
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
